# user/newsboat.timer
[Unit]
Description=Newsboat RSS Feed Update Daemon Timer

[Timer]
OnActiveSec=3h
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
